# lliurex-analytics-server/etc/cron.weekly/analytics_restart
#!/bin/bash
systemctl stop lliurex-analytics-server
sleep 20
systemctl start lliurex-analytics-server
