# lliurex-analytics-server/etc/cron.weekly/analytics_restart
#!/bin/bash
service lliurex-analytics-server stop
sleep 20
service lliurex-analytics-server start